# scripts/macos/com.optionsrisk.koreareport.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  Korea pre-open status report: weekdays at 08:30 LOCAL time.
  Assumes the Mac mini timezone is KST. Replace /Users/REPLACE_ME with your path.
-->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.optionsrisk.koreareport</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/REPLACE_ME/options-risk-alert/scripts/report.sh</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/REPLACE_ME/options-risk-alert</string>

  <key>StartCalendarInterval</key>
  <array>
    <dict><key>Weekday</key><integer>1</integer><key>Hour</key><integer>8</integer><key>Minute</key><integer>30</integer></dict>
    <dict><key>Weekday</key><integer>2</integer><key>Hour</key><integer>8</integer><key>Minute</key><integer>30</integer></dict>
    <dict><key>Weekday</key><integer>3</integer><key>Hour</key><integer>8</integer><key>Minute</key><integer>30</integer></dict>
    <dict><key>Weekday</key><integer>4</integer><key>Hour</key><integer>8</integer><key>Minute</key><integer>30</integer></dict>
    <dict><key>Weekday</key><integer>5</integer><key>Hour</key><integer>8</integer><key>Minute</key><integer>30</integer></dict>
  </array>

  <key>StandardOutPath</key>
  <string>/Users/REPLACE_ME/options-risk-alert/data/report.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/REPLACE_ME/options-risk-alert/data/report.err.log</string>
</dict>
</plist>
